# ansible/roles/gns3_host/files/com.user.loginscript.plist
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.user.loginscript</string>
    <key>ProgramArguments</key>
    <array><string>/Users/user/.vboxnet0.zsh</string></array>
    <key>RunAtLoad</key>
    <true/>
  </dict>
</plist>
